# tarsnap.target
[Unit]
Description=Tarsnap target
StopWhenUnneeded=yes
